# hostApp/services/pytender-kim.service
[Unit]
Description=Pytender API
After=network.target

[Service]
ExecStart=/usr/bin/python -u keepItMoving.py
WorkingDirectory=/usr/local/lib/PytTender/orders
StandardOutput=inherit
StandardError=inherit
Restart=always
User=pytend

[Install]
WantedBy=multi-user.target
